# test/supervisord.conf
[supervisord]
;nodaemon=true

[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
autorestart=true

[program:redis]
command=/usr/bin/redis-server
autorestart=true
